# .infrastructure/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log                ; supervisord log file
logfile_maxbytes=5MB                            ; maximum size of logfile before rotation
logfile_backups=2                               ; number of backed up logfiles
loglevel=info                                   ; info, debug, warn, trace
nodaemon=true                                   ; run supervisord as a daemon
user=root                                       ; default user
childlogdir=/var/log                            ; where child log files will live

[program:cron]
command=/entrypoint.cron.sh
user=root
stdout_events_enabled=true
stderr_events_enabled=true

[program:infracheck]
command=/entrypoint.cmd.sh
stdout_events_enabled=true
stderr_events_enabled=true
